# gs_pipeline/docker/supervisord.conf
; gs_pipeline: run the Streamlit UI and the trainer watcher side-by-side.
; supervisord here owns both processes; if either dies, supervisord restarts
; it. The container exits only if supervisord itself stops.

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:ui]
command=streamlit run /app/gs_pipeline/ui/app.py --server.maxUploadSize 20480
environment=PYTHONPATH=/app
priority=20
autostart=true
autorestart=true
stdout_logfile=/data/logs/ui.out
stderr_logfile=/data/logs/ui.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
stopsignal=TERM
stopwaitsecs=10

[program:worker]
command=python -m gs_pipeline.trainer.watcher
    --inbox /data/inbox
    --work /data/work
    --outbox /data/outbox
    --logs /data/logs
    --config-yaml /app/gs_pipeline/trainer/config.yaml
    --quality-preset Auto
environment=PYTHONPATH=/app
priority=30
autostart=true
autorestart=true
stdout_logfile=/data/logs/worker.out
stderr_logfile=/data/logs/worker.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
stopsignal=TERM
stopwaitsecs=30
